# _CONTROL/launchdaemons/com.familiendash.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.familiendash</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Shared/GIT/familiendash/.venv/bin/python</string>
        <string>/Users/Shared/GIT/familiendash/scripts/feedback_api.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Shared/GIT/familiendash</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/Shared/GIT/familiendash/logs/api.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/GIT/familiendash/logs/api_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/vpn-flodido/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/vpn-flodido</string>
    </dict>
</dict>
</plist>
